# debian/trackkr-webapp.postinst
#!/bin/sh

/usr/sbin/useradd --system --home=/usr/share/trackkr --shell=/bin/nologin trackkr

ln -s /etc/nginx/sites-available/trackkr.conf /etc/nginx/sites-enabled/

chown -R trackkr:trackkr /var/log/trackkr

invoke-rc.d trackkr.webapp restart
invoke-rc.d nginx restart
