# debian/trackkr-daemon.postinst
#!/bin/sh

invoke-rc.d trackkr.trackkrd restart
